# deploy/gsim2.service
[Unit]
Description=GSIM2 — TessTOC Goldratt simulator
Documentation=https://github.com/svtriers-dot/GSIM2
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=gsim
WorkingDirectory=/opt/gsim2
EnvironmentFile=/opt/gsim2/.env
ExecStart=/usr/bin/node /opt/gsim2/dist/index.cjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gsim2

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gsim2

[Install]
WantedBy=multi-user.target
